# usr/lib/systemd/system/refind-btrfs.timer
[Unit]
Description=Generate rEFInd manual boot stanzas from Btrfs snapshots

#[Timer]
#OnBootSec=1min

[Install]
WantedBy=timers.target
