# etc/systemd/system/yadm_commit_and_push.timer
[Unit]
Description=Atualiza repositório Dotfiles(https://github.com/rizzini/Dotfiles)

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
